# autossh/autossh-ultra.service
[Unit]
Description=Persistent SSH tunnel to ultra.cc services
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=debian
Environment=AUTOSSH_GATETIME=0
EnvironmentFile=/etc/default/dashboard-autossh
ExecStart=/usr/bin/autossh -M 0 -N -o ExitOnForwardFailure=yes -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -o StrictHostKeyChecking=yes -o UserKnownHostsFile=/home/debian/.ssh/known_hosts -i ${AUTOSSH_IDENTITY_FILE} -L 127.0.0.1:16141:127.0.0.1:16141 -L 127.0.0.1:16124:127.0.0.1:16124 ${AUTOSSH_TARGET}
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
